# docker-builds/storm/storm-supervisor/sshd.supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/data/storm/%(program_name)s.out.log
stderr_logfile=/data/storm/%(program_name)s.err.log
autorestart=true
